# files/0-system/usr/lib/udev/rules.d/zz1-usb-storage-optimized-async-udev.rules
ACTION=="add|change", SUBSYSTEM=="block", SUBSYSTEMS=="usb", RUN+="/usr/bin/usb-dirty-pages-udev '%k' '%E{ID_USB_VENDOR_ID}' '%E{ID_USB_MODEL_ID}'"
