# files/systemd/system/adw-gtk3-installer-privileged.timer
[Unit]
Description=Timer to ensure that adw-gtk3 theme is applied for system flatpaks

[Timer]
OnBootSec=20

[Install]
WantedBy=timers.target
