# config/systemd/system/nightly-reboot.service
[Unit]
Description=Scheduled Reboot

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --force reboot
