# config/systemd/user/flatpak-gidro@deluge.service.d/deluge.conf
[Service]
Environment="FLATPAK_RUN_OPTIONS=org.deluge_torrent.deluge"
